# crontab
* * * * * php /var/www/bin/console about